# passworddsync/passworddsync.service
[Unit]
Description=Simple Password synchronisation daemon
After=network-pingable.service passwordd.service network-online.target
Wants=network-pingable.service passwordd.service

[Service]
ExecStart=/usr/bin/passworddsync-script start
ExecStop=/usr/bin/passworddsync-script stop

[Install]
WantedBy=multi-user.target
