# passworddsync/passworddsync.service
[Unit]
Description=Simple Password synchronisation daemon
After=network.target passwordd.service

[Service]
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/passworddsync-script start
ExecStop=/usr/bin/passworddsync-script stop
ExecStartPost=/bin/sleep 5

[Install]
WantedBy=multi-user.target
